# config/messenger-worker.ini
[program:messenger-consume]
command=php74-cli /home/serwer2255151/public_html/HSBG/bin/console messenger:consume async  --limit=5 --memory-limit=200M
user=serwer2255151
startsecs=0
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d